# dockerfiles/docker-irisnet/supervisor/conf.d/supervisor-iris.conf
[program:iris]
command=iris start --home %(ENV_IRIS_HOME)s
redirect_stderr=true