# dockerfiles/docker-kava/supervisor/conf.d/supervisor-kvd.conf
[program:kvd]
command=kvd start --home %(ENV_KVD_HOME)s
redirect_stderr=true