# dockerfiles/docker-regen/supervisor/conf.d/supervisor-xrnd.conf
[program:xrnd]
command=xrnd start --home %(ENV_REGEN_HOME)s
redirect_stderr=false
autostart=true
autorestart=unexpected